# server/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

# ── FalkorDB ──────────────────────────────────────────────────────
[program:falkordb]
command=/usr/local/bin/redis-server --loadmodule /opt/falkordb/falkordb.so --dir /data --appendonly yes --save 60 1 --protected-mode no --bind 0.0.0.0 --port 6379
autostart=true
autorestart=true
priority=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── Memory Service (doit démarrer après FalkorDB) ─────────────────
[program:memory]
command=python3 -m uvicorn main:app --host 0.0.0.0 --port 5003
directory=/app/memory
autostart=true
autorestart=true
priority=2
startsecs=3
environment=FALKORDB_HOST="127.0.0.1",FALKORDB_PORT="6379"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── Agent Orya ────────────────────────────────────────────────────
[program:agent-orya]
command=python3 -m uvicorn main:app --host 0.0.0.0 --port 5001
directory=/app/agent-orya
autostart=true
autorestart=true
priority=3
startsecs=2
environment=MEMORY_URL="http://127.0.0.1:5003"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── Orchestrator ──────────────────────────────────────────────────
[program:orchestrator]
command=python3 -m uvicorn main:app --host 0.0.0.0 --port 5002
directory=/app/orchestrator
autostart=true
autorestart=true
priority=4
startsecs=2
environment=MEMORY_URL="http://127.0.0.1:5003",GATEWAY_URL="http://127.0.0.1:4001"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── Gateway (Hono/Bun) ───────────────────────────────────────────
[program:gateway]
command=/root/.bun/bin/bun src/index.ts
directory=/app/gateway
autostart=true
autorestart=true
priority=5
startsecs=2
environment=AGENT_ORYA_URL="http://127.0.0.1:5001",ORCHESTRATOR_URL="http://127.0.0.1:5002",GATEWAY_PORT="4001"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
