# server/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

# ── Neo4j (priority 1) ────────────────────────────────────────────
[program:neo4j]
command=/usr/share/neo4j/bin/neo4j console
autostart=true
autorestart=true
priority=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── PostgreSQL (priority 1) ───────────────────────────────────────
[program:postgres]
command=/usr/local/bin/init-postgres.sh
autostart=true
autorestart=true
priority=1
user=root
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── Graphiti REST server (priority 2) ─────────────────────────────
[program:graphiti-server]
command=bash -c "python3 /app/wait-for-neo4j.py && exec python3 -m uvicorn main:app --host 0.0.0.0 --port 8000"
directory=/app/graphiti-server
autostart=true
autorestart=true
priority=2
startsecs=5
environment=PYTHONUNBUFFERED="1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── Agent (LangGraph + FastAPI) (priority 3) ──────────────────────
[program:agent]
command=bash -c "python3 /app/wait-for-neo4j.py && exec python3 -m uvicorn agent.main:app --host 0.0.0.0 --port 5001"
directory=/app
autostart=true
autorestart=true
priority=3
startsecs=5
environment=PYTHONUNBUFFERED="1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# ── Gateway (Hono/Bun WebSocket) (priority 4) ─────────────────────
[program:gateway]
command=/root/.bun/bin/bun src/index.ts
directory=/app/gateway
autostart=true
autorestart=true
priority=4
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
